# vidin-calafat-bridge.crontab
0 0 0 1/1 * ? * cd /home/suricactus/work/geo/vidin-calafat-bridge; bash ./vidin-calafat-bridge.sh